# ansible/roles/backend/templates/gunicorn.service.j2
[Unit]
Description={{ app_name }} backend
After=network.target

[Service]
Type=notify
User={{ deploy_user }}
Group={{ deploy_group }}
WorkingDirectory={{ backend_dir }}
ExecStart={{ backend_dir }}/venv/bin/gunicorn \
          --workers {{ backend_workers }} \
          --worker-class sync \
          --worker-connections 1000 \
          --max-requests 1000 \
          --max-requests-jitter 50 \
          --timeout 30 \
          --bind 127.0.0.1:{{ backend_port }} \
          --access-logfile /var/log/{{ app_name }}/gunicorn-access.log \
          --error-logfile /var/log/{{ app_name }}/gunicorn-error.log \
          --log-level info \
          recruitment_backend.wsgi:application

ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=on-failure
RestartSec=5

# Security settings
NoNewPrivileges=true
LimitNOFILE=65536

# Environment
Environment="DJANGO_SETTINGS_MODULE={{ django_settings_module }}"
EnvironmentFile={{ backend_dir }}/.env

[Install]
WantedBy=multi-user.target